# config/includes.chroot/etc/cron.d/start_csrd_cloud
@reboot root /opt/csrd-cloud/csrd.sh >> /var/log/csrd.log 2>&1
